# supervisord.conf
[supervisord]

nodaemon=true

[program:win-app]

command=bee run
